# systemd/edge-motion.service
[Unit]
Description=Edge Motion (touchpad edge scrolling helper)
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/edge-motion --threshold 0.06 --hold-ms 80 --pulse-ms 10 --pulse-step 3
Restart=always
RestartSec=2
KillSignal=SIGTERM

# This daemon reads touchpad events and emits mouse events via uinput,
# so it intentionally runs with elevated privileges by default.
User=root

[Install]
WantedBy=multi-user.target
