# systemd/edge-motion.service
[Unit]
Description=Edge Motion (touchpad edge scrolling helper)
After=multi-user.target

[Service]
Type=simple
EnvironmentFile=-/etc/default/edge-motion
ExecStartPre=@BINDIR@/edge-motion-auto-update
ExecStart=@BINDIR@/edge-motion $EDGE_MOTION_ARGS
Restart=always
RestartSec=2
KillSignal=SIGTERM

# This daemon reads touchpad events and emits mouse events via uinput,
# so it intentionally runs with elevated privileges by default.
User=root

[Install]
WantedBy=multi-user.target
